# docker/s6-overlay/s6-rc.d/laravel-scheduler/run
#!/command/with-contenv sh

cd /var/www/html
exec s6-setuidgid www-data php artisan schedule:work
